# docker/container/supervisor/common/conf.d/200-php-fpm.conf
[program:php-fpm]
priority=200
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/sbin/php-fpm -R
numprocs=1
numprocs_start=1
startsecs=0
autostart=true
autorestart=true
stopasgroup=true
stopwaitsecs=30
killasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
