# docker/container/supervisor/common/conf.d/400-worker.conf
[program:worker]
priority=400
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /var/www/html/artisan queue:work database --daemon --tries=3 --timeout=90 --sleep=3 --max-time=3600 --backoff=15
user=www-data
numprocs=0
numprocs_start=1
startsecs=0
autostart=true
autorestart=true
stopasgroup=true
stopwaitsecs=30
killasgroup=true
redirect_stderr=false
stderr_logfile=/var/log/supervisor/worker.log
stderr_logfile_maxbytes=0
stdout_logfile=/var/log/supervisor/worker.log
stdout_logfile_maxbytes=0
